# etc/systemd/netmon@.service
[Unit]
Description=NetMon Per-Connection Server

[Service]
ExecStart=-/usr/local/bin/netmon-agent
StandardInput=socket
